# jellylink.service
[Unit]
Description=JellyLink Media Organizer
After=network.target

[Service]
Type=simple
User=martin
Group=martin
WorkingDirectory=/home/martin/jellylink
ExecStart=/usr/bin/python3 /home/martin/jellylink/jellylink.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jellylink

[Install]
WantedBy=multi-user.target
